# deployment/dnsmasq-gui.service
[Unit]
Description=DNSmasq GUI Web Interface
After=network.target
Wants=network.target

[Service]
Type=simple
User=dnsmasq-gui
Group=dnsmasq-gui
WorkingDirectory=/opt/dnsmasq-gui
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=-/opt/dnsmasq-gui/.env

# DNSmasq configuration file paths
Environment=DNSMASQ_CONFIG_PATH=/etc/dnsmasq.conf
Environment=DNSMASQ_HOSTS_FILE=/etc/dnsmasq.hosts
Environment=DNSMASQ_LEASES_PATH=/var/lib/dhcp/dhcpd.leases
Environment=DNSMASQ_STATIC_LEASES_CONFIG_FILE=/etc/dnsmasq.d/dnsmasq-static-leases.conf
Environment=DNSMASQ_RANGES_CONFIG_FILE=/etc/dnsmasq.d/dnsmasq-ranges.conf
Environment=DNSMASQ_OPTIONS_CONFIG_FILE=/etc/dnsmasq.d/dnsmasq-options.conf
Environment=DNSMASQ_ADVANCED_CONFIG_FILE=/etc/dnsmasq.d/dnsmasq-advanced.conf

# Security settings (relaxed for file access)
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/etc/dnsmasq.conf /etc/dnsmasq.hosts /etc/dnsmasq.d /var/lib/dhcp /var/log

[Install]
WantedBy=multi-user.target
